# installers/ubuntu/24.04/statix-agent.service
[Unit]
Description=Statix Agent
Documentation=https://statix.pettersson.online/docs/nodes/manual-installation
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment=STATIX_AGENT_CONFIG=/etc/statix/agent.json
Environment=STATIX_VERSION_FILE=/opt/statix/version.json
Environment=STATIX_APPLY_WIREGUARD=0
Environment=XDG_RUNTIME_DIR=/run/statix-agent
EnvironmentFile=-/etc/statix/agent.env
ExecStart=/usr/local/bin/statix run
Restart=always
RestartSec=5s
KillSignal=SIGINT
TimeoutStopSec=20s

User=statix-agent
Group=statix-agent

# LXC jobs run as unprivileged containers, but LXC still needs namespace
# creation, writable delegated cgroups, and the setuid uidmap helpers for the
# configured subordinate IDs.
NoNewPrivileges=false
Delegate=yes
PrivateTmp=true
PrivateDevices=false
ProtectHome=read-only
ProtectSystem=strict
# Ubuntu's LXC package uses this multiarch package directory as a temporary
# rootfs mountpoint while starting containers. Keep it writable even though the
# rest of /usr is read-only in the service mount namespace.
ReadWritePaths=-/usr/lib/x86_64-linux-gnu/lxc/rootfs -/usr/lib/aarch64-linux-gnu/lxc/rootfs
ProtectControlGroups=false
# LXC must mount procfs/sysfs for job containers. Making kernel tunables
# read-only in the agent's own mount namespace can make the container's procfs
# mount look less restrictive and fail with EPERM.
ProtectKernelTunables=false
ProtectKernelModules=true
ProtectClock=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=false
RestrictRealtime=true
RestrictNamespaces=false
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
ReadOnlyPaths=/etc/statix /opt/statix
RuntimeDirectory=statix-agent
RuntimeDirectoryMode=0700
StateDirectory=statix-agent

[Install]
WantedBy=multi-user.target
